# etc/profile.d/runindent_options_kr.sh
# Puts a list of indent options for runindent into the environment
#
# Dot this file to make runindent format according to the style of the well-
# known book "The C Programming Language" by Kernighan & Ritchie
#
# To make this file the installation default, symlink to it from
# /etc/profile.d, or copy it there.
#
RUNINDENT_OPTIONS='-kr'
export RUNINDENT_OPTIONS
